# rootfs/etc/services.d/crond/run
#!/usr/bin/execlineb -P

crond -f -L /var/log/cron -c /etc/cron.d
